# deploy/systemd/bingwall-api.service
[Unit]
Description=BingWall FastAPI service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bingwall
Group=bingwall
SupplementaryGroups=www-data
WorkingDirectory=/opt/bingwall/app
EnvironmentFile=/etc/bingwall/bingwall.env
ExecStart=/opt/bingwall/app/.venv/bin/python -m uvicorn app.main:create_app --factory --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/bingwall /var/log/bingwall /etc/bingwall

[Install]
WantedBy=multi-user.target
