# deploy/systemd/bingwall-api.service
[Unit]
Description=BingWall FastAPI service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bingwall
Group=bingwall
SupplementaryGroups=www-data
WorkingDirectory=/opt/bingwall/app
EnvironmentFile=/etc/bingwall/bingwall.env
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/env uv run --no-sync python -m uvicorn app.main:create_app --factory --host ${BINGWALL_APP_HOST} --port ${BINGWALL_APP_PORT}
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30
UMask=0027
NoNewPrivileges=true
RemoveIPC=true
PrivateTmp=true
PrivateDevices=true
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectProc=invisible
ProcSubset=pid
ProtectSystem=strict
ReadWritePaths=/var/lib/bingwall /var/log/bingwall /etc/bingwall
CapabilityBoundingSet=
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
