# deploy/systemd/bingwall-api.service
[Unit]
Description=BingWall FastAPI service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bingwall
Group=bingwall
SupplementaryGroups=www-data
WorkingDirectory=/opt/bingwall/app
EnvironmentFile=/etc/bingwall/bingwall.env
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/env uv run --no-sync python -m uvicorn app.main:create_app --factory --host ${BINGWALL_APP_HOST} --port ${BINGWALL_APP_PORT}
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/bingwall /var/log/bingwall /etc/bingwall

[Install]
WantedBy=multi-user.target
